# templates/usr/lib/systemd/system/couchdb-backup.service.j2
[Unit]
Description=CouchDB Backeruper

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/couchdb-backup.sh
		